# srcs/docker/influxdb/rootfs/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true
user=root
 
[program:influxd]
command=/usr/sbin/influxd
autorestart=true

[program:telegraf]
command=/usr/bin/telegraf
autorestart=true
