# srcs/docker/grafana/rootfs/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true
user=root
 
[program:grafana]
command=grafana-7.3.3/bin/grafana-server
autorestart=true

[program:telegraf]
command=/usr/bin/telegraf
autorestart=true
